# sblite/root/etc/hotplug.d/iface/99-sblite
#!/bin/sh

[ "$ACTION" = "ifup" -o "$ACTION" = "ifupdate" ] && [ $(uci get sblite.main.enable) = "1" ] && {
    # if "$DEVICE" in outbounds restart sblite
}